# supervisord-haproxy.conf
[program:haproxy]
command=/etc/init.d/haproxy start
numprocs=1
autostart=true
autorestart=false
