# supervisord-haproxy.conf
[program:haproxy]
command=service haproxy start
numprocs=1
autostart=true
autorestart=true
